# etc/uci-defaults/99-default-settings
#!/bin/sh

nvram set CountryCode=EU
nvram commit

exit 0
