# uci-defaults/99-default-settings
#!/bin/sh

nvram set CountryCode=CN
nvram commit

uci set wireless.wifi0.country='US'
uci set wireless.wifi1.country='US'
uci set wireless.wifi2.country='US'
uci commit

exit 0
