# uci-defaults/99-default-settings
#!/bin/sh

nvram set CountryCode=CN
nvram commit

exit 0
